# init/smack.mount
[Unit]
Description=Smack filesystem mounting
ConditionSecurity=smack

[Mount]
What=smackfs
Where=/smack
Type=smackfs
